# invitation-app/invitation-public.service
[Unit]
Description=Invitation App - Public RSVP Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/invitation-app
ExecStart=/home/pi/invitation-app/venv/bin/python /home/pi/invitation-app/public_server.py
Restart=always
RestartSec=5
Environment=PATH=/home/pi/invitation-app/venv/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
